# systemd/gunicorn002.socket
[Unit]
Description=gunicorn002 socket

[Socket]
ListenStream=/run/gunicorn002.sock

[Install]
WantedBy=sockets.target

